# supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/bin/bash -c "exec /usr/sbin/sshd -D"

[program:freeswitch]
command=/bin/bash -c "exec /opt/freeswitch/bin/freeswitch -ncwait"


