# supervisord.conf
[supervisord]
nodaemon=true


[program:freeswitch]
command=/bin/bash -c "exec /opt/freeswitch/bin/freeswitch -ncwait -D FOREGROUND"
